# polybot.service
[Unit]
Description=Polymarket Bot Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aadichauhan
WorkingDirectory=/home/aadichauhan/polymarket_bot
ExecStart=/home/aadichauhan/polymarket_bot/venv/bin/python daemon.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
